# lib/chromeos/arc/container/myfiles/arc-myfiles.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Mount /home/chronos/user/MyFiles for ARC."
author        "chromium-os-dev@chromium.org"

start on started arcpp-post-login-services
stop on stopped arcpp-post-login-services

env SOURCE=/home/chronos/user/MyFiles
env DEST=/run/arc/media/MyFiles
env UMASK=007
env ANDROID_MEDIA_UID=1023
env ANDROID_MEDIA_GID=1023

# syslog-cat is used to redirect stdio from sdcard to the logging daemon.
# minijail is used in the mount-passthrough-jailed script.
exec /usr/sbin/syslog-cat --identifier="${UPSTART_JOB}" -- \
  mount-passthrough-jailed "${SOURCE}" "${DEST}" "${UMASK}" \
    "${ANDROID_MEDIA_UID}" "${ANDROID_MEDIA_GID}" "full"

post-stop exec logger -t "${UPSTART_JOB}" "Post-stop ${UPSTART_JOB}"
