# lib/chromeos/arc/setup/init/arc-boot-continue.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Turn the container into a fully functional one"
author        "chromium-os-dev@chromium.org"

# This is to block arcpp-post-login-services until this task finishes.
start on starting arcpp-post-login-services
stop on stop-arc-instance or stopping ui

# This job needs to be a task to block session_manager
task

# The following environment variables are passed from session_manager
import CHROMEOS_DEV_MODE
import CHROMEOS_INSIDE_VM
import CHROMEOS_USER
import CONTAINER_PID
import COPY_PACKAGES_CACHE
import DEMO_SESSION_APPS_PATH
import DISABLE_BOOT_COMPLETED_BROADCAST
import ENABLE_ADB_SIDELOAD
import ENABLE_ARC_NEARBY_SHARE
import IS_DEMO_SESSION
import LOCALE
import MANAGEMENT_TRANSITION
import PREFERRED_LANGUAGES
import SKIP_PACKAGES_CACHE_SETUP
import SKIP_GMS_CORE_CACHE_SETUP

exec /sbin/minijail0 -c \
    cap_chown,cap_dac_read_search,cap_dac_override,cap_fowner,cap_sys_admin+e \
    -- /usr/sbin/arc-setup --mode=boot-continue "--log_tag=${UPSTART_JOB}"
