# lib/chromeos/arc/setup/init/arcpp-post-login-services.conf
# Copyright 2021 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Job to start ARC++ services after user login"
author        "chromium-os-dev@chromium.org"

start on continue-arc-boot
stop on stop-arc-instance or stopping ui

# The following environment variables are passed from session_manager
# and are imported from the event that starts the job.
import CHROMEOS_DEV_MODE
import CHROMEOS_INSIDE_VM
import CHROMEOS_USER
import CONTAINER_PID
import COPY_PACKAGES_CACHE
import DEMO_SESSION_APPS_PATH
import DISABLE_BOOT_COMPLETED_BROADCAST
import ENABLE_ADB_SIDELOAD
import ENABLE_ARC_NEARBY_SHARE
import IS_DEMO_SESSION
import LOCALE
import MANAGEMENT_TRANSITION
import PREFERRED_LANGUAGES
import SKIP_PACKAGES_CACHE_SETUP
import SKIP_GMS_CORE_CACHE_SETUP

export CHROMEOS_DEV_MODE
export CHROMEOS_INSIDE_VM
export CHROMEOS_USER
export CONTAINER_PID
export COPY_PACKAGES_CACHE
export DEMO_SESSION_APPS_PATH
export DISABLE_BOOT_COMPLETED_BROADCAST
export ENABLE_ADB_SIDELOAD
export ENABLE_ARC_NEARBY_SHARE
export IS_DEMO_SESSION
export LOCALE
export MANAGEMENT_TRANSITION
export PREFERRED_LANGUAGES
export SKIP_PACKAGES_CACHE_SETUP
export SKIP_GMS_CORE_CACHE_SETUP
