# lib/chromeos/shill/init/shill-start-user-session.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Perform connection-manager-specific session setup tasks."
author          "chromium-os-dev@chromium.org"

start on load-shill-profile
# NB: session manager does not wait for this task to complete.
task

# Sacrifice before OOM panic. Default score for platform services that
# can recover.
oom score -100

# CHROMEOS_USER is an upstart-provided variable setting passed from
# session_manager via the 'load-shill-profile' signal.
import CHROMEOS_USER
exec shill_login_user ${CHROMEOS_USER}
