# lib/chromeos/shill/init/shill_respawn.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Restart shill at slow rate if it fails to respawn"
author        "chromium-os-dev@chromium.org"

# If shill fails to start up due to fast-respawn, this is a bad
# thing.  Connectivity in chromium-os is absolutely essential, so
# respond to this failure by pausing and restarting.
start on stopped shill PROCESS=respawn

# Sacrifice before OOM panic. Default score for platform services that
# can recover.
oom score -100

script
  sleep 60
  start shill
end script
