# lib/chromeos/shill/init/shill.service.in
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

[Unit]
Description=Run the shill network connection manager
Before=network.target
After=@expected_started_services@
Requisite=@expected_started_services@
Conflicts=pre-shutdown.service
[Service]
# Default parameters to be overridden by systemd configuration.
Environment="SHILL_LOG_LEVEL=0" \
	"SHILL_LOG_SCOPES=" \
	"SHILL_PASSIVE_MODE=" \
	"SHILL_ACCEPT_HOSTNAME_FROM=" \
	"SHILL_MINIMUM_MTU=" \
	"BLOCKED_DEVICES=" \
	"SHILL_TEST_DEVICES="
Restart=on-failure

ExecStartPre=/bin/sh /usr/share/cros/init/shill-pre-start.sh

ExecStart=/bin/sh /usr/share/cros/init/shill.sh

ExecStopPost=/usr/sbin/bootstat shill-stop
