# lib/chromeos/shill/init/shill-stop-user-session.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Perform connection-manager-related session cleanup tasks."
author          "chromium-os-dev@chromium.org"

# Use 'start on stopping' and 'task' to hold up restarting the session
# manager until the connection manager has cleaned up
# session-dependent state.
start on stopping ui
task

# Sacrifice before OOM panic. Default score for platform services that
# can recover.
oom score -100

script
  # In case there are quick, back-to-back login/logout -- logout should win.
  stop shill-start-user-session || :
  exec shill_logout_user
end script
