# lib/chromeos/arc/apk-cache/init/apk-cache-cleaner.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Performs cleaning of the APK cache directory"
author          "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping system-services

exec periodic_scheduler --period=86400 --timeout=600 \
    --task_name=apk-cache-cleaner-jailed -- /usr/sbin/apk-cache-cleaner-jailed
