# lib/chromeos/arc/vm/scripts/init/mount-media-dirs/arcvm-mount-sdcard-dir.conf
# Copyright 2021 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

#TODO(b/177367496): Remove this script once /run/arcvm is merged to /run/arc.

description    "Mount /run/arc/sdcard as a shared mount point in ARCVM."
author         "chromium-os-dev@chromium.org"

# Adjust the lifetime to that of seneschal to make "Play files" shareable with
# Linux. Note that it covers the whole lifetime of arcvm-mount-play-files.
start on starting seneschal
stop on stopped seneschal

# Do not allow OOM killer to terminate this job. This job should be started and
# terminated with the seneschal service.
oom score never

env SDCARD_DIR=/run/arc/sdcard

pre-start script
  # Create the runtime directory for Android's "Play files" accessed by
  # seneschal and the Files app. Make it a shared mount point so that mount
  # events inside mount-passthrough's namespace are propagated to the init
  # namespace. We do this here because the directory needs to be mounted
  # before seneschal starts.
  # UID and GID of the directory are set to android-root:android-everybody.
  # TODO(b/148298890): Use arc-setup instead of mere shell commmands.
  mkdir -p "${SDCARD_DIR}"
  mount -o rw,nodev,noexec,nosuid,mode=0755,uid=655360,gid=665357 \
        -t tmpfs tmpfs "${SDCARD_DIR}"
  mount --make-shared "${SDCARD_DIR}"
end script

post-stop script
  # TODO(b/148298890): Use arc-setup instead of a mere shell commmand.
  umount -R "${SDCARD_DIR}" || true
end script
