# lib/chromeos/arc/vm/scripts/init/mount-media-dirs/arcvm-mount-myfiles.conf
# Copyright 2020 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description    "Mount the MyFiles directory using mount-passthrough."
author         "chromium-os-dev@chromium.org"

start on starting arcvm-post-login-services
stop on stopping arcvm-post-login-services

# Do not allow OOM killer to terminate this job since it will break file
# sharing support (which e.g. our Photos app relies on) until ARCVM restarts.
# TODO(b/164939771): Choose a better OOM score which guarantees the Chrome
# browser process is always terminated before this one.
oom score never

# VmHWM of the process is ~6 MB in the steady state, and setting the VM size
# soft limit below that results in process termination. On the other hand,
# VmPeak points to somewhat higher value (~200 MB), and if we specify a soft
# limit lower than that, while not causing process termination, it results in
# memory allocation errors for actions like `ls` command.
# Given that, we choose to set a high value of 1 GB as the soft limit here.
# TODO(b/164936284): Specify a lower limit after improving mount-passthrough.
limit as 1000000000 unlimited

import CHROMEOS_USER

env DEST=/run/arcvm/media/MyFiles
env UMASK=007
env CONCIERGE_NAMESPACE_FILE=/run/namespaces/mnt_concierge

script
  if [ -z "${CHROMEOS_USER}" ]; then
    logger -t "${UPSTART_JOB}" "No Chrome OS user is specified"
    exit 1
  fi

  IS_MOUNTED=$(cryptohome --action=is_mounted --user="${CHROMEOS_USER}")
  if [ "${IS_MOUNTED}" != "true" ]; then
    logger -t "${UPSTART_JOB}" "Cryptohome is not mounted for the user"
    exit 1
  fi

  CRYPTOHOME_USER_HOME=$(cryptohome-path user "${CHROMEOS_USER}")
  # TODO(b/178353348): Avoid using the hard-coded MyFiles path.
  SOURCE="${CRYPTOHOME_USER_HOME}/MyFiles"

  if [ ! -d "${SOURCE}" ]; then
    logger -t "${UPSTART_JOB}" \
      "MyFiles directory ${SOURCE} does not exist"
    exit 1
  fi

  # Create the destination directory.
  nsenter --mount="${CONCIERGE_NAMESPACE_FILE}" --no-fork \
    -- mkdir -p "${DEST}"

  exec /usr/sbin/syslog-cat --identifier="${UPSTART_JOB}" -- \
       mount-passthrough-jailed-media "${SOURCE}" "${DEST}" "${UMASK}"
end script

post-stop script
  # ${DEST} is not automatically unmounted when terminated by SIGKILL.
  nsenter --mount="${CONCIERGE_NAMESPACE_FILE}" --no-fork \
    -- umount "${DEST}" || true
end script
