# lib/chromeos/sirenia/upstart/dugong.conf
# Copyright 2021 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start dugong, the sirenia chrome os daemon"
author          "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services
respawn
respawn limit 3 10  # if the job respawns 3 times in 10 seconds, stop trying.

# Run dugong with the default vsock address and port
exec minijail0 -u dugong -- /usr/bin/dugong -U vsock://vsock:2:5552
